# tools/deployments/gunicorn/etc/systemd/system/mkdevs.socket
[Unit]
Description=mkdevs socket

[Socket]
ListenStream=/run/mkdevs.sock

[Install]
WantedBy=sockets.target
